# systemd-units/btrfs-snap@frequent.timer
[Unit]
PartOf=btrfs-snap.target
Persistent=true

[Timer]
OnCalendar=*:0/15
