# systemd-units/btrfs-snap@.timer
[Unit]
PartOf=btrfs-snap.target

[Timer]
Persistent=true
OnCalendar=%i
