# systemd-units/btrfs-snap@frequent.timer
[Unit]
PartOf=btrfs-snap.target

[Timer]
Persistent=true
OnCalendar=*:0/15
